# roles/common/files/systemd/pacman-download.service
[Unit]
Description=Download Pacman packages

[Service]
ExecStart=/usr/bin/pacman -Syuw --noconfirm
