# roles/zfs/templates/syncoid.service.j2
[Unit]
Description=Run syncoid

[Service]
ExecStart=/usr/bin/syncoid --recursive --no-sync-snap --sendoptions="w" --mbuffer-size=1G {{ item.src }} {{ item.dest }}
