# roles/zfs/templates/syncoid.timer.j2
[Unit]
Description=Run syncoid {{ item.frequency }}

[Timer]
{{ item.timer }}

[Install]
WantedBy=timers.target
