# roles/dotfiles/files/.config/autostart/devilspie.desktop
[Desktop Entry]
Exec=/usr/bin/devilspie
Icon=
Name=devilspie
Path=
Terminal=False
Type=Application
